# systemd/radio-kiosk.service
[Unit]
Description=Radio Kiosk Backend Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/radio-kiosk/src/service.py
Restart=always
RestartSec=5
WorkingDirectory=/opt/radio-kiosk
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
